# pipeline/deploy/daily-sync.service
[Unit]
Description=SRRL daily sync — upload images + meteo to HuggingFace
Documentation=https://huggingface.co/datasets/knl2366/NREL_Sky_Imagery
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=exouser
WorkingDirectory=/home/exouser/srrl-pipeline/resources
ExecStart=/usr/bin/python3 /home/exouser/srrl-pipeline/resources/daily_sync.py
TimeoutStartSec=1800

# Environment
Environment=SRRL_IMAGE_DIR=/media/volume/Primary-Dataset
Environment=SRRL_METEO_DIR=/media/volume/Primary-Dataset/meteorological
EnvironmentFile=-/home/exouser/.env

# Logging
StandardOutput=append:/var/log/srrl/daily-sync.log
StandardError=append:/var/log/srrl/daily-sync.log

# Hardening
NoNewPrivileges=true
